# etc/s6-overlay/s6-rc.d/sshd/run
#!/command/with-contenv bash
set -euo pipefail
exec /usr/sbin/sshd -D -e
