# etc/s6-overlay/s6-rc.d/kasmvnc/run
#!/command/with-contenv bash
set -euo pipefail
exec /bin/bash /root/setup_kasmvnc.sh
