# etc/s6-overlay/s6-rc.d/socat/run
#!/command/with-contenv bash
set -euo pipefail
exec /usr/bin/socat TCP-LISTEN:9223,reuseaddr,fork TCP:localhost:9222,nonblock
